# vnstatd-rs.service
[Unit]
Description=vnStatd-rs network traffic monitor
After=network.target

[Service]
Type=forking
RuntimeDirectory=vnstatd-rs
PIDFile=/run/vnstatd-rs/vnstatd-rs.pid
User=vnstat
Group=vnstat
WorkingDirectory=/var/lib/vnstat-rs 
ExecStart=/usr/bin/vnstatd-rs -d --pidfile /run/vnstatd-rs/vnstatd-rs.pid
Restart=always

[Install]
WantedBy=multi-user.target
